# scripts/remote_runner_service_migration/systemd/meridian-alpha-mainnet-unattended-daily-policy.service
# Draft only. Install disabled in a future approved remote apply window.
[Unit]
Description=Meridian Alpha daily unattended policy gate for one bounded same-slot two-phase mainnet sequence
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/meridian_alpha_live_runner/repo
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/meridian_alpha_live_runner/bin/with-live-env /root/meridian_alpha_live_runner/venv/bin/python scripts/live_trading/run_hv_balanced_mainnet_unattended_daily_policy.py --config config/live_trading/hv_balanced_binance_usdm_live_supervisor_multiphase_wallet2_reserve200_meridian_remote_runner.yaml --apply --enable-supervisor-timer --run-health-monitor
TimeoutStartSec=1200
StandardOutput=journal
StandardError=journal
