# scripts/remote_runner_service_fix_window/systemd-dropins/meridian-alpha-mainnet-health-monitor.service.d/10-meridian-path.conf
# Meridian path/config-resolution fix. Install only in a separate approved fix window.
[Service]
WorkingDirectory=/root/meridian_alpha_live_runner/repo
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONNOUSERSITE=1
ExecStart=
ExecStart=/root/meridian_alpha_live_runner/bin/with-live-env /usr/bin/env PYTHONPATH=/root/meridian_alpha_live_runner/repo/src PYTHONNOUSERSITE=1 VIRTUAL_ENV=/root/meridian_alpha_live_runner/venv PATH=/root/meridian_alpha_live_runner/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /root/meridian_alpha_live_runner/venv/bin/python /root/meridian_alpha_live_runner/repo/scripts/live_trading/run_hv_balanced_mainnet_health_monitor.py --config /root/meridian_alpha_live_runner/repo/config/live_trading/hv_balanced_binance_usdm_live_supervisor_multiphase_wallet2_reserve200_meridian_handoff_observation.yaml --systemd-timer-name meridian-alpha-mainnet-supervisor-live.timer
